# bwall.service
[Unit]
Description=Bing Wallpaper Downloader

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bwall today