# bwall.timer
[Unit]
Description=Run Bing wallpaper downloader daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target